# barcode-scanner.service
[Unit]
Description=Commercial Barcode Scanner Web Service
Documentation=https://github.com/your-repo/barcode-scanner
After=network.target
Wants=network.target

[Service]
Type=simple
User=gt39112
Group=gt39112
WorkingDirectory=/var/www/html/abhimanyu/barcode_scanner_clean
Environment=PATH=/var/www/html/abhimanyu/barcode_scanner_clean/venv/bin
Environment=PYTHONPATH=/var/www/html/abhimanyu/barcode_scanner_clean/src
ExecStart=/var/www/html/abhimanyu/barcode_scanner_clean/venv/bin/python web_app.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
TimeoutStopSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/html/abhimanyu/barcode_scanner_clean /var/log/barcode-scanner

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=barcode-scanner

[Install]
WantedBy=multi-user.target
