# mqtt-pi-client.service
[Unit]
Description=MQTT Pi Client - Device Announcement Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/barcode_scanner_clean/src/utils
ExecStart=/usr/bin/python3 /home/pi/barcode_scanner_clean/src/utils/mqtt_pi_client.py iot.caleffionline.it
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/pi/barcode_scanner_clean/src

[Install]
WantedBy=multi-user.target
