# mqtt-discovery.service
[Unit]
Description=MQTT Device Discovery Service
After=network.target mosquitto.service
Wants=network.target
Requires=mosquitto.service

[Service]
Type=simple
User=azureuser
Group=azureuser
WorkingDirectory=/var/www/html/abhimanyu/barcode_scanner_clean
ExecStart=/usr/bin/python3 /var/www/html/abhimanyu/barcode_scanner_clean/src/utils/mqtt_device_discovery.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

Environment=PYTHONPATH=/var/www/html/abhimanyu/barcode_scanner_clean/src

[Install]
WantedBy=multi-user.target
